# test-monitoring.service
[Unit]
Description=Monitor the 'test' process and report to monitoring server
After=network.target

[Service]
Type=simple
ExecStart=/path/to/script.sh
Restart=on-failure
User=root
Group=root
Environment=PATH=/usr/bin:/usr/local/bin
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
